# packer/target/logtrixia/setup/logtrixia.service
[Unit]
Description=Logtrixia 
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/etc/logtrixia
ExecStart=/usr/local/bin/docker-compose -f elastic-docker-tls.yml up -d
ExecStop=/usr/local/bin/docker-compose -f elastic-docker-tls.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target